# policycoreutils/mcstrans/src/mcstrans.service
[Unit]
Description=Translates SELinux MCS/MLS labels to human readable form
ConditionSecurity=selinux

[Service]
ExecStart=/sbin/mcstransd -f

[Install]
WantedBy=multi-user.target
